# systemd/fsbackup-logrotate-metric.timer
[Unit]
Description=fsbackup logrotate metric timer

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target

