# systemd/fsbackup-runner-daily@.timer
[Unit]
Description=fsbackup daily runner timer (%i)

[Timer]
# Default schedule — overridden per-class by fs-schedule-apply.sh
# drop-ins in /etc/systemd/system/fsbackup-runner-daily@<class>.timer.d/schedule.conf
OnCalendar=*-*-* 02:00
RandomizedDelaySec=5m
Persistent=true

[Install]
WantedBy=timers.target
