# systemd/fsbackup-runner-weekly@.timer
[Unit]
Description=fsbackup weekly runner timer (%i)

[Timer]
# Default schedule — overridden per-class by fs-schedule-apply.sh
OnCalendar=Sat *-*-* 03:00
RandomizedDelaySec=5m
Persistent=true

[Install]
WantedBy=timers.target
