# systemd/fsbackup-mirror-daily.service
[Unit]
Description=fsbackup mirror (daily pass) to SNAPSHOT_MIRROR_ROOT
After=fsbackup-runner@class1.service fsbackup-runner@class2.service
Before=fsbackup-retention.service

[Service]
UMask=0022
Type=oneshot
User=fsbackup
Group=fsbackup
ExecStart=/opt/fsbackup/bin/fs-mirror.sh daily

