# systemd/fsbackup-promote.timer
[Unit]
Description=fsbackup promotion check

[Timer]
OnCalendar=03:30
Persistent=true

[Install]
WantedBy=timers.target

