# systemd/fsbackup-annual-promote.timer
[Unit]
Description=fsbackup annual promote timer

[Timer]
OnCalendar=*-01-05 03:00:00
Persistent=true

[Install]
WantedBy=timers.target


