# systemd/fsbackup-logrotate-metric.service
[Unit]
Description=fsbackup logrotate health metric

[Service]
UMask=0022
Type=oneshot
User=root
ExecStart=/opt/fsbackup/bin/fs-logrotate-metric.sh
