# systemd/fsbackup-mirror-promote.service
[Unit]
Description=fsbackup mirror (post-promote pass) to SNAPSHOT_MIRROR_ROOT
After=fsbackup-promote.service

[Service]
UMask=0022
Type=oneshot
User=fsbackup
Group=fsbackup
ExecStart=/opt/fsbackup/bin/fs-mirror.sh promote

