# systemd/fsbackup-runner@class3.timer
[Unit]
Description=fsbackup monthly snapshot timer (class3)

[Timer]
OnCalendar=*-*-01 04:45
RandomizedDelaySec=10m
Persistent=true

[Install]
WantedBy=timers.target
